# deployment/systemd/pdf-tools.service
[Unit]
Description=PDF Tools FastAPI Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=pdftools
Group=pdftools

WorkingDirectory=/opt/pdf-tools
EnvironmentFile=/etc/pdf-tools/pdf-tools.env

ExecStart=/opt/pdf-tools/venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --proxy-headers \
    --forwarded-allow-ips=192.168.3.0/24 \
    --workers 2 \
    --no-access-log

Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

MemoryDenyWriteExecute=false
UMask=0027

[Install]
WantedBy=multi-user.target
